# futurehome/rootfs/etc/services.d/futurehome/finish
#!/usr/bin/with-contenv bashio
# Nothing special – just allow S6 to stop gracefully
exit 0
